# rootfs/etc/services.d/postgres/run
#!/command/with-contenv bash
set -euo pipefail

case "${KHOJ_USE_INTERNAL_POSTGRES:-true}" in
    false|FALSE|False|0|no|NO|No)
        exec sleep infinity
        ;;
esac

if [ -n "${POSTGRES_HOST:-}" ]; then
    exec sleep infinity
fi

PGDATA="/var/lib/postgresql/data"
PG_VERSION="$(ls /usr/lib/postgresql | sort -V | tail -n1)"

if [[ -f "${PGDATA}/PG_VERSION" ]]; then
    DATA_PG_VERSION="$(head -n1 "${PGDATA}/PG_VERSION" | tr -d '[:space:]')"
    if [[ -n "${DATA_PG_VERSION}" ]]; then
        if [[ -x "/usr/lib/postgresql/${DATA_PG_VERSION}/bin/postgres" ]]; then
            PG_VERSION="${DATA_PG_VERSION}"
        else
            echo "[khoj-aio] Existing PGDATA requires PostgreSQL ${DATA_PG_VERSION}, but no compatible server is installed." >&2
            exit 1
        fi
    fi
fi

exec su postgres -s /bin/sh -c "/usr/lib/postgresql/${PG_VERSION}/bin/postgres -D \"$PGDATA\" -c listen_addresses='127.0.0.1'"
